# ssh_ids_notify.service
[Unit]
Description=Send Message To Mattermost

[Service]
User=ibcadmin
ExecStart=/<PATH TO>/send_message.sh

[Install]
WantedBy=multi-user.target

